# server_config/www2/init/celery.conf
start on runlevel [2345]
stop on runlevel [!2345]

exec su -c '/home/mhc/mhc-venv/bin/python /home/mhc/mhc-backend/manage.py celeryd --loglevel=info -f /home/mhc/logs/celeryd.log -s /home/mhc/celery-schedules -B -v 2 -E --concurrency=8 --maxtasksperchild=1' - mhc

respawn
